# deploy/com.mp3fbf.twitter-processor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Twitter Bookmark Processor - Daemon Mode (Polling)

  This plist runs the processor in daemon mode, polling data/backlog/
  every 2 minutes for new Twillot export files.

  Installation:
    1. Run ./setup-macos.sh to install (auto-configures paths)
    2. Or manually: chmod +x deploy/run-processor.sh && cp deploy/*.plist ~/Library/LaunchAgents/
    3. launchctl load ~/Library/LaunchAgents/com.mp3fbf.twitter-processor.plist

  API Key:
    The wrapper script fetches ANTHROPIC_API_KEY from 1Password at runtime.
    Requires: op CLI authenticated (service account or biometric)

  Management:
    launchctl list | grep twitter        # Check if running
    launchctl stop com.mp3fbf.twitter-processor   # Stop daemon
    launchctl start com.mp3fbf.twitter-processor  # Start daemon
    launchctl unload ~/Library/LaunchAgents/com.mp3fbf.twitter-processor.plist  # Remove

  Logs:
    tail -f /tmp/twitter-processor.log   # stdout
    tail -f /tmp/twitter-processor.err   # stderr
-->
<plist version="1.0">
<dict>
    <!-- Unique service identifier -->
    <key>Label</key>
    <string>com.mp3fbf.twitter-processor</string>

    <!-- Program to execute (wrapper script fetches API key from 1Password) -->
    <key>ProgramArguments</key>
    <array>
        <string>/workspace/twitter-bookmark-processor/deploy/run-processor.sh</string>
        <string>--insight</string>
        <string>--source</string>
        <string>both</string>
    </array>

    <!-- Working directory (where src/ is located) -->
    <key>WorkingDirectory</key>
    <!-- TODO: Update to your project directory -->
    <string>/workspace/twitter-bookmark-processor</string>

    <!-- Environment variables -->
    <!-- Note: Wrapper script loads OP token from ~/.secrets/op-token-yolo -->
    <!-- Note: API keys are fetched from 1Password at runtime -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>TWITTER_OUTPUT_DIR</key>
        <string>__PROJECT_DIR__/../notes/Sources/twitter/</string>

        <key>LOG_LEVEL</key>
        <string>INFO</string>
    </dict>

    <!-- Run once at load, then every 15 minutes -->
    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>900</integer>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/tmp/twitter-processor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/twitter-processor.err</string>

    <!-- Nice level (lower priority, be gentle on system) -->
    <key>Nice</key>
    <integer>5</integer>

    <!-- Process type (background = low priority) -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
