# deploy/com.mp3fbf.twitter-processor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Twitter Bookmark Processor - Insight Engine Daemon

  Polls X API + Twillot backlog every 15 min, processes through
  Insight Engine (capture + Opus distill), writes Obsidian notes.

  Installation:
    bash deploy/setup-launchd.sh

  Logs:
    tail -f /tmp/twitter-processor.log
    tail -f /tmp/twitter-processor.err
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mp3fbf.twitter-processor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/workspace/twitter-bookmark-processor/deploy/run-processor.sh</string>
        <string>--insight</string>
        <string>--source</string>
        <string>both</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/workspace/twitter-bookmark-processor</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TWITTER_OUTPUT_DIR</key>
        <string>__PROJECT_DIR__/../notes/Sources/twitter/</string>

        <key>LOG_LEVEL</key>
        <string>INFO</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/twitter-processor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/twitter-processor.err</string>

    <key>Nice</key>
    <integer>5</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
